# templates/upstart.conf.j2
#!upstart

# This file was generated by Ansible for {{ ansible_fqdn }}
# Do NOT modify this file by hand!

description "Start celery workers for {{celery_app_name}}"

start on (local-filesystems and runlevel [2345])
stop on runlevel [!2345]

kill timeout 20


{% if item.user|default(celery_user) == 'root' %}
env C_FORCE_ROOT=true
env HOME=/root
{% endif %}

env CNAME="{{item.action|default('worker')}}{{item.queue|default('') and ('-' + item.queue) or ''}}"
env QUEUE="{{ '-Q %s' % item.queue if item.queue|default(None) else '' }}"
env APP="{{ '-A %s' % item.app_module|default(celery_app_module) if item.app_module|default(celery_app_module) else '' }}"
env CONCURRENCY="{{ '-c %s' % item.concurrency|default(celery_concurrency) if item.action == 'worker' else '' }}"

setuid {{ item.user|default(celery_user) }}

respawn

limit nofile 65536 65536
console log
chdir {{item.app_dir|default(celery_app_dir)}}

exec {{celery_bin}} {{item.action}} $APP $QUEUE $CONCURRENCY \
    -f {{item.log_dir|default(celery_log_dir)}}/$CNAME.log -l {{item.log_level|default(celery_log_level)}}

